# root/etc/s6-overlay/s6-rc.d/init-secrets/up
#!/usr/bin/with-contenv bash

echo "Creating secrets environment file..."

# Create environment file that will be sourced by other services
SECRETS_ENV="/tmp/secrets.env"

# Process legacy Docker Secrets (backwards compatibility)
if [[ -d /run/secrets ]]; then
    echo "Processing Docker Secrets..."

    if [[ -f /run/secrets/tandoor_secret_key ]]; then
        SECRET_KEY=$(cat /run/secrets/tandoor_secret_key)
        echo "export SECRET_KEY='$SECRET_KEY'" >> "$SECRETS_ENV"
        echo "✓ Set SECRET_KEY from Docker secrets"
    fi

    if [[ -f /run/secrets/tandoor_postgres_password ]]; then
        POSTGRES_PASSWORD=$(cat /run/secrets/tandoor_postgres_password)
        echo "export POSTGRES_PASSWORD='$POSTGRES_PASSWORD'" >> "$SECRETS_ENV"
        echo "✓ Set POSTGRES_PASSWORD from Docker secrets"
    fi

    if [[ -f /run/secrets/tandoor_postgres_user ]]; then
        POSTGRES_USER=$(cat /run/secrets/tandoor_postgres_user)
        echo "export POSTGRES_USER='$POSTGRES_USER'" >> "$SECRETS_ENV"
        echo "✓ Set POSTGRES_USER from Docker secrets"
    fi

    echo "Docker secrets processed."
fi

# Make environment file readable by all services
chmod 644 "$SECRETS_ENV"

echo "✓ Secrets environment file created at $SECRETS_ENV"