# root/etc/s6-overlay/s6-rc.d/init-branding/up
#!/usr/bin/with-contenv bash

# Display Mildman1848 ASCII art branding
cat /etc/s6-overlay/s6-rc.d/init-adduser/branding
